# oxide.service
[Unit]
Description=Oxide Racing Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/data/Code/oxide
Environment=NODE_ENV=production
# Uncomment and edit the next line to load environment variables from a file
# EnvironmentFile=/data/Code/oxide/.env
ExecStart=/usr/bin/node /data/Code/oxide/dist/server.js
Restart=on-failure
RestartSec=5
User=www-data
# Adjust the user as needed

[Install]
WantedBy=multi-user.target

# Service: oxide

# This service runs the C++ multi-port server (`my-socket-server`).
#
# - Build with CMake as described in the README.
# - The service listens on ports 3000 (HTTP), 8226/8228/7003 (Custom1), and 43300 (Custom2).
#
# ---
#
# **Note:** This service file is for the C++/CMake version of the project.
